# documentacion/script.ini
[program:script]
command=/bin/bash -c "sleep 10 && /bin/php /var/www/html/script.php"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/script.err.log
stdout_logfile=/var/log/supervisor/script.out.log


[program:script2]
command=/bin/bash -c "sleep 10 && /bin/php /var/www/html/script2.php"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/script.errScript2.log
stdout_logfile=/var/log/supervisor/script.outScript2.log